# centreon/files/etc/inittab
# /etc/inittab

::sysinit:/sbin/openrc sysinit
::wait:/sbin/openrc boot
::wait:/sbin/openrc default

